# master/kube-controller-manager.service
[Service]
ExecStart=/home/core/share/demo/bin/kube-controller-manager --master=127.0.0.1:8080
